# SOURCES/mwlib.ini
[group:mwlib]
programs  = nserve,nslave,postman,qserve

[program:nserve]
command         = mw-nserve
user            = nobody
directory       = /var/cache/mwlib
redirect_stderr = true
stdout_logfile  = /var/log/mwlib/nserve.log
[program:nslave]
command         = mw-nslave --cachedir /var/cache/mwlib
user            = nobody
directory       = /var/cache/mwlib
redirect_stderr = true
stdout_logfile  = /var/log/mwlib/nslave.log
[program:postman]
command         = mw-postman
user            = nobody
directory       = /var/cache/mwlib
redirect_stderr = true
stdout_logfile  = /var/log/mwlib/postman.log
[program:qserve]
command         = mw-qserve
user            = nobody
directory       = /var/cache/mwlib
redirect_stderr = true
stdout_logfile  = /var/log/mwlib/qserve.log